# ansible/templates/stratosentry.service.j2
[Unit]
# StratoSentry NIDS sensor — managed by Ansible. Do not edit by hand.
Description=StratoSentry NIDS Sensor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Runs as root because scapy needs raw-socket access for packet sniffing.
User=root
# Network interface to sniff, auto-detected by the playbook for this host.
Environment=STRATOSENTRY_IFACE={{ nids_interface }}
WorkingDirectory={{ project_dir }}
# Use the virtualenv's interpreter so dependencies resolve correctly.
ExecStart={{ venv_dir }}/bin/python {{ project_dir }}/src/sensor.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
